# sonos_upstart.conf
# Sonos.js upstart script
# http://caolanmcmahon.com/posts/deploying_node_js_with_upstart

description "Sonos monitor"
author      "glen"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn

exec sudo -u pi sh -c "/usr/local/bin/node /home/pi/sonos/sonos.js"